# autostart/ibus-daemon.desktop
[Desktop Entry]
Exec=/usr/bin/ibus-daemon -drx
Name=ibus-daemon
Type=Application
