# autostart/unclutter.desktop
[Desktop Entry]
Exec=/usr/bin/unclutter -idle 5
Name=unclutter
Comment=hide idle mouse pointer
Type=Application
